# systemd/infra-maintenance.service
[Unit]
Description=Infrastructure Maintenance Task

[Service]
Type=oneshot
ExecStart=/root/linux-infra-intern-assignment/scripts/maintenance.sh
